# supervisord.conf
[supervisord]
nodaemon=true

[program:jupyter]
command=jupyter lab --ip=0.0.0.0 --no-browser --port 8888 --NotebookApp.token=${JUPYTER_TOKEN}

[program:rstudio]
command=/init
